# codeworkers/crontab
35 07  * * * root /scrapy-crawl.sh >> /var/log/scrapy-cron.log 2 >&1
